# contrib/discourse-server.conf

description "Discourse Server"


start on discourse

chdir /home/discourse/discourse

env RAILS_ENV=production
exec bundle exec rails s -p 3001 >> logs/thin.log
