# contrib/discourse-server.conf

description "Discourse Server"


start on starting discourse
stop on stopping discourse

chdir /home/discourse/discourse

env RAILS_ENV=production
exec bundle exec rails s -p 3001 >> logs/thin.log
